# alpha-wrap/share/chroot_master-builders/files/setup
#!/sbin/openrc-run
# Auto-installation service, to be run at the beginning to perform setup.
# If installation is successful, script is not called again.

# Setup configuration BEGIN

# <==- Edition specific variable declarations place holder, do NOT remove NOR edit this commend -==>

# END OF Setup configuration

# read kernel options
init_KOPT() {
	eval "set -- $(cat /proc/cmdline 2>/dev/null)"
	for opt; do
		case "$opt" in
			ssh_*=*)
				eval "KOPT_${opt%%=*}='${opt#*=}'" ;;
		esac
	done
}

start() {

    ls -aFl /boot

    echo $DEVD
	rm -f /etc/runlevels/*/$RC_SVCNAME
	init_KOPT
	local rc=0
	ebegin "Starting ${RC_SVCNAME}"
	if [ -n "$KOPT_ssh_key" ] && [ ! -f "/root/.ssh/authorized_keys" ]; then
		einfo "Fetching ssh keys"
		mkdir -pm 700 /root/.ssh
		checkpath -fm 0600 /root/.ssh/authorized_keys
		case "$KOPT_ssh_key" in
			https://*|ftps://*|http://*)
				wget -q "$KOPT_ssh_key" -O /root/.ssh/authorized_keys
				rc=$?;;
			*) echo "$KOPT_ssh_key" > /root/.ssh/authorized_keys;;
		esac
	fi

	#
	eend $rc
}

