# files/autoinstall
#!/sbin/openrc-run

depend() {
        need localmount sysfs cgroups
        after lxd firewall ntp-client
}




start() {
        # _control 'start' 'RUNNING'
        einfo "Please login as 'admin' (no password)"
        einfo "    this will finalise setup."

        einfo ""
        einfo "If you did setup ssh-key you can login"
        einfo "    ssh -i key_path admin@102.303.342.094"

        eend 0
}

stop() {
        # _control 'stop' 'STOPPED'
        eend 0
}

